# bigip_exporter.service
[Unit]
Description=Prometheus BIG-IP Exporter
Documentation=https://github.com/Haameed/f5_bigip_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bigip_exporter
Group=bigip_exporter
ExecStart=/usr/local/bin/f5_bigip_exporter \
    -config /etc/f5_bigip_exporter/bigip-config.yaml \
    -listen :11000 \
    -scrape-timeout 30
Restart=on-failure
RestartSec=5s

# --- Hardening ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
ReadOnlyPaths=/etc/bigip_exporter

[Install]
WantedBy=multi-user.target
